# etc/init.d/awesant-agent.in
#!/bin/sh

### BEGIN INIT INF
# Provides:         awesant
# Required-Start:   $all
# Required-Stop:    $all
# Default-Start:    2 3 4 5
# Default-Stop:     0 1 6
# Short-Description: Starts awesant agent
# chkconfig: 2345 80 15
# Description: Awesant services
### END INIT INFO

BASENAME=$(basename $0)
CONFIG=@@CONFIG@@
RUNDIR=@@RUNDIR@@
PIDFILE=$RUNDIR/agent.pid
DAEMON=@@DAEMON@@

[ -f /etc/sysconfig/$BASENAME ] && . /etc/sysconfig/$BASENAME
[ -f /etc/default/$BASENAME ] && . /etc/default/$BASENAME

if test ! -d "$RUNDIR" ; then
    mkdir -p $RUNDIR
fi

__check_awesant_agent() {
    DAEMON_EXIST=""
    if test -e "$PIDFILE" ; then
        DAEMON_PID=$(cat $PIDFILE)
        if test -n "$DAEMON_PID" ; then
            DAEMON_EXIST=$(ps $DAEMON_PID | grep awesant)
        fi
    fi
}

__start() {
    echo -n "Starting ${BASENAME}: "
    if test -z "$DAEMON_EXIST" ; then
        $DAEMON --config "$CONFIG" --pidfile "$PIDFILE" &

        for i in `seq 1 3` ; do
            __check_awesant_agent
            if test ! -z "$DAEMON_EXIST" ; then
                break
            fi
            sleep 1
        done

        if test -z "$DAEMON_EXIST" ; then
            echo "awesant failed to start"
            RETVAL=1
        else
            echo "awesant started successfully, PID $DAEMON_PID"
            RETVAL=0
        fi
    else
        echo "awesant-agent already running"
        RETVAL=0
    fi
}

__stop() {
    echo -n "Stopping ${BASENAME}: "
    if test -z "$DAEMON_EXIST" ; then
        echo "awesant-agent not running"
        RETVAL=0
    else
        kill -15 $DAEMON_PID

        for i in `seq 1 15` ; do
            __check_awesant_agent
            if test -z "$DAEMON_EXIST" ; then
                break
            fi
            sleep 1
        done

        if test -z "$DAEMON_EXIST" ; then
            echo "awesant-agent stopped successfully"
            RETVAL=0
            rm -f ${PIDFILE}
        else
            echo "ERR: awesant-agent was not stopped successfully, PID $DAEMON_PID"
            RETVAL=1
        fi
    fi
    return $RETVAL
}

__status() {
    echo -n "Status ${BASENAME}: "
    if test -z "$DAEMON_EXIST" ; then
        echo "awesant-agent is not running"
        RETVAL=1
    else
        echo "awesant-agent is running, PID $DAEMON_PID"
        RETVAL=0
    fi
}

__usage() {
    echo "Usage: $BASENAME {start|stop|restart}"
    RETVAL=1
}

__check_awesant_agent

case "$1" in
    start   ) __start ;;
    stop    ) __stop ;;
    status  ) __status ;;
    restart ) __stop && __start ;;
    *       ) __usage ;;
esac

exit $RETVAL
